# DEBIAN/auto-cpulimit.service
[Unit]
Description="auto-cpulimit package"
