# DEBIAN/auto-cpulimit.service
[Unit]
Description=auto-cpulimit package

[Service]
ExecStart=/usr/bin/acl.sh

[Install]
WantedBy=multi-user.target
